# deploy/systemd/meetyou-desktop-client.service.template
[Unit]
Description=MeetYou Desktop Client
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{MEETYOU_USER}}
Group={{MEETYOU_GROUP}}
WorkingDirectory={{WORKDIR}}
EnvironmentFile={{ENV_FILE}}
Environment=PYTHONUNBUFFERED=1
ExecStart={{PYTHON_BIN}} -m desktop_client
Restart=always
RestartSec=3
TimeoutStopSec=30
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=false
ReadWritePaths={{WORKDIR}} {{WORKDIR}}/user {{WORKDIR}}/logs
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
